# networkd-state@.target
[Unit]
Description=Interface state target